# crates/rulesteward-auditd/tests/corpus/auditd/rocky9-deep-many-watches/audit.rules
## This file is automatically generated from /etc/audit/rules.d
# Deep / many watches (catalog #39): dozens of watches on long, deep paths.
# Parser + ranking stress (distinct from #16's mixed huge ruleset: this is
# watch-only, with deep nested paths). Each is a single-file write/attr watch
# (900 B/event, f3 section 3.3). One watch is a RECURSIVE deep DIRECTORY watch
# (HIGH, recursive subtree, f3 section 2.2) which must out-rank the file watches.
-w /etc/security/limits.d/99-custom-nofile.conf -p wa -k deepwatch_0
-w /etc/systemd/system/multi-user.target.wants/custom-app.service -p wa -k deepwatch_1
-w /usr/lib/python3.11/site-packages/app/config/settings.py -p wa -k deepwatch_2
-w /var/lib/pgsql/16/data/postgresql.conf -p wa -k deepwatch_3
-w /opt/vendor/product/etc/conf.d/security/policy.d/main.policy -p wa -k deepwatch_4
-w /etc/pki/tls/private/server-internal-ca-chain.key -p wa -k deepwatch_5
-w /etc/NetworkManager/system-connections/corp-vpn.nmconnection -p wa -k deepwatch_6
-w /etc/audit/rules.d/10-base-config-very-long-name.rules -p wa -k deepwatch_7
-w /etc/ssh/sshd_config.d/50-redhat-hardening-overrides.conf -p wa -k deepwatch_8
-w /usr/local/share/ca-certificates/internal/root-ca-2026.crt -p wa -k deepwatch_9
-w /etc/sssd/conf.d/01-ldap-domain-corp-example-com.conf -p wa -k deepwatch_10
-w /etc/firewalld/zones/internal-restricted-zone.xml -p wa -k deepwatch_11
-w /etc/selinux/targeted/policy/policy.33 -p wa -k deepwatch_12
-w /etc/containers/registries.conf.d/000-shortnames-internal.conf -p wa -k deepwatch_13
-w /etc/dnf/plugins/subscription-manager-internal-mirror.conf -p wa -k deepwatch_14
-w /etc/grub.d/40_custom-secure-boot-entries -p wa -k deepwatch_15
-w /etc/cron.d/internal-compliance-scan-nightly -p wa -k deepwatch_16
-w /etc/logrotate.d/custom-application-audit-logs -p wa -k deepwatch_17
-w /etc/modprobe.d/blacklist-legacy-network-drivers.conf -p wa -k deepwatch_18
-w /etc/sysctl.d/99-kernel-hardening-overrides.conf -p wa -k deepwatch_19
# A RECURSIVE deep directory watch (HIGH; out-ranks the single-file watches).
-w /etc/security/ -p wa -k security_tree
