# crates/rulesteward-auditd/tests/corpus/auditd/rocky8-pci-dss/audit.rules
## This file is automatically generated from /etc/audit/rules.d
-w /etc/passwd -p wa -k pci_identity
-w /etc/shadow -p wa -k pci_identity
-w /bin/ -p wa -k pci_binaries
-w /sbin/ -p wa -k pci_binaries
-w /etc/cron.d/ -p wa -k pci_jobs
